# lxc/ej-livestream-origin/etc/cron.d/ej_hls_cleanup
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* *	* * *	root	/root/scripts/livestream_cleanup.sh >/dev/null 2>&1
