# host/etc/cron.d/ej_update
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

11 3    * * *   root    /root/ej_scripts/update_debian.sh >/dev/null 2>&1
11 4    * * *   root    /root/ej_scripts/update_container.sh >/dev/null 2>&1
